# ansible/roles/kubernetes-master/templates/kube-scheduler.conf.j2
description "kube-scheduler service"

start on runlevel [2345]
stop on runlevel [!2345]

kill timeout 60

exec kube-scheduler \
	--address=0.0.0.0 \
	--kubeconfig=/srv/kubernetes/kubeconfig \
	--v=1
